# mtc-keymaster/server2/c/mtc-ca.service
[Unit]
Description=MTC CA/Log Server (C)
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
ExecStart=/usr/local/bin/mtc_server \
    --port 8444 \
    --data-dir /home/ubuntu/.mtc-ca-data \
    --tokenpath /home/ubuntu/.env \
    --ca-name MTC-CA-C \
    --log-id 32473.2 \
    --tls-cert /home/ubuntu/.mtc-ca-data/server-cert.pem \
    --tls-key /home/ubuntu/.mtc-ca-data/server-key.pem \
    --ech-name factsorlie.com \
    --log-level 4 \
    --dh-port 8445 \
    --mqc-port 8446 \
    --tpm-path /home/ubuntu/.TPM/factsorlie.com-ca
Restart=on-failure
RestartSec=5
Environment=HOME=/home/ubuntu
Environment=LD_LIBRARY_PATH=/usr/local/lib
WorkingDirectory=/home/ubuntu/postWolf/mtc-keymaster/server2/c

[Install]
WantedBy=multi-user.target
